# scripts/com.worldcup.autocommit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent that runs scripts/auto-commit.sh every 30 minutes to push vault + model data
  to GitHub. Install: cp this file to ~/Library/LaunchAgents/ and
    launchctl load -w ~/Library/LaunchAgents/com.worldcup.autocommit.plist
  Uninstall:
    launchctl unload -w ~/Library/LaunchAgents/com.worldcup.autocommit.plist
  Logs: /tmp/worldcup-autocommit.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.worldcup.autocommit</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/cashmcdearis/Documents/worldcup-moneylines/scripts/auto-commit.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/worldcup-autocommit.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/worldcup-autocommit.log</string>
</dict>
</plist>
